# ibus/.config/systemd/user/ibus.service
[Unit]
Description=IBus Daemon

[Service]
ExecStart=/usr/bin/ibus-daemon -rx

[Install]
WantedBy=default.target
