# ibus/.config/systemd/user/openrgb.service
[Unit]
Description=OpenRGB

[Service]
ExecStart=/usr/bin/openrgb --profile default

[Install]
WantedBy=default.target
